# launchd/com.chetan.kharcha-mini-ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chetan.kharcha-mini-ingest</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/bun</string>
    <string>run</string>
    <string>/Users/chetanjain/Desktop/code/kharcha-mini/src/ingest/run.ts</string>
  </array>
  <key>StartInterval</key>
  <integer>900</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>AbandonProcessGroup</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/chetanjain/Library/Logs/kharcha-mini-ingest.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/chetanjain/Library/Logs/kharcha-mini-ingest.err.log</string>
</dict>
</plist>
